# robot/server/etc/rc.local
python /opt/pi/robot/server/server.py &
