# robot/server/etc/rc.local
LD_LIBRARY_PATH=/opt/pi/robot/server/mjpg-streamer/mjpg-streamer python /opt/pi/robot/server/server.py &
